# trust-collapse/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log

[program:coordinator]
command=python3 /opt/services/coordinator/server.py
autostart=true
autorestart=true
stdout_logfile=/var/log/coordinator.log
stderr_logfile=/var/log/coordinator-error.log

[program:builder]
command=python3 /opt/services/builder/server.py
autostart=true
autorestart=true
stdout_logfile=/var/log/builder.log
stderr_logfile=/var/log/builder-error.log
